# tts/02_rk3588_sherpa_onnx_tts_deploy_project/systemd/sherpa-onnx-tts.service
[Unit]
Description=RK3588 sherpa-onnx TTS Service
After=network-online.target sound.target
Wants=network-online.target
[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
Environment=SHERPA_TTS_ENV={{PROJECT_ROOT}}/config/sherpa_tts.env
ExecStart={{PROJECT_ROOT}}/scripts/05_start_api.sh
Restart=always
RestartSec=3
[Install]
WantedBy=multi-user.target
